# deploy/systemd/arenatech-resolve-stale-conversations.timer
[Unit]
Description=Dispara resolve-stale-conversations de hora em hora

[Timer]
OnCalendar=*:00
AccuracySec=30s
# Persistent: se o VPS ficou fora no horario, roda assim que voltar.
Persistent=false

[Install]
WantedBy=timers.target
